# 03-measurement-infra/templates/systemd/cloudflare-measurement.timer
[Unit]
Description=zdns-collectors

[Timer]
OnCalendar=*-*-* *:*:0
Persistent=true
AccuracySec=15

[Install]
WantedBy=timers.target
  
